# packages/eventlircd/config/irexec.service
[Unit]
Description=Eventlircd daemon

[Service]
ExecStart=/usr/bin/irexec /etc/lirc/lircrc

[Install]
WantedBy=multi-user.target
